# agent/examples/pmeow-agent.service
[Unit]
Description=PMEOW Agent Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pmeow-agent
EnvironmentFile=/etc/pmeow-agent/pmeow-agent.env
ExecStart=/opt/pmeow-agent/.venv/bin/pmeow-agent run
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
